# files/upstart_uwsgi.conf
description "uWSGI Emperor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env UWSGI=/usr/local/bin/uwsgi
env LOGTO=/var/log/uwsgi/emperor.log

exec $UWSGI --master --emperor /etc/uwsgi/apps-enabled --logto $LOGTO --die-on-term --pidfile /tmp/uwsgi-emperor.pid
